# pulseaudio.service.d/20-droid-wait.conf
[Service]
# Wait for miniaf to start up...
ExecStartPre=/usr/bin/waitforservice init.svc.miniaf
