# pulseaudio.service.d/20-droid-wait.conf
[Service]
# Wait for miniaf to start up...
ExecStartPre=/bin/sh -c '/bin/test -e /system/bin/miniafservice && /usr/bin/waitforservice init.svc.miniaf || true'
